# core/deploy/backup.timer
# Run deploy/backup.sh every four hours. Enable with: systemctl enable --now backup.timer (setup.sh does this).
# Persistent=true so a missed calendar run fires after the box returns instead of silently extending the RPO.
[Unit]
Description=Four-hour backup of the nullsink billing DBs

[Timer]
OnCalendar=*-*-* 00/4:00:00 UTC
RandomizedDelaySec=15min
Persistent=true

[Install]
WantedBy=timers.target
