# core/deploy/backup.service
# One-shot DB backup, run daily by backup.timer. Snapshots balances.db + pending.db via sqlite3 .backup,
# optionally age-encrypts, and (if BACKUP_PUSH_CMD is set) ships the artifact off-box. Runs AS THE SERVICE
# USER so the SQLite sidecars it touches stay service-owned (a root .backup breaks billing writes).
[Unit]
Description=nullsink billing-DB backup (one-shot; run by backup.timer)
# Page Telegram if a backup run fails (a silent backup failure = creeping data-loss risk).
OnFailure=status-alert@%n.service
# Never snapshot the short-lived half-restored pair left by a killed/power-lost restore. restore.sh stops the
# timer too; this persistent condition is the reboot backstop until an operator safely resumes the restore.
ConditionPathExists=!/var/lib/nullsink/.restore-in-progress
# During post-restore health activation the pair is valid, but restore still holds the exclusive ledger lock.
# A Persistent timer may fire immediately when restarted; skip that transient launch successfully (no
# OnFailure page). restore.sh removes this marker only after readiness, releases the lock, then explicitly
# starts one fresh backup so the skipped calendar occurrence is not lost.
ConditionPathExists=!/var/lib/nullsink/.restore-activation-pending
ConditionPathExists=!/var/lib/nullsink/.deploy-in-progress

[Service]
Type=oneshot
User=nullsink
Group=nullsink
WorkingDirectory=/opt/nullsink
# BACKUP_AGE_RECIPIENT / BACKUP_PUSH_CMD / BACKUP_KEEP live here; `-` so a missing file is fine.
EnvironmentFile=-/etc/nullsink.env
# Same StateDirectory as both app services, so /var/lib/nullsink (and backups/ under it) is writable under
# ProtectSystem=strict. balances.db belongs to nullsink-proxy and pending.db to nullsink-payments, but they
# share one uid and one directory, so a single backup run still snapshots both consistently.
StateDirectory=nullsink
StateDirectoryMode=0700
UMask=0077
Environment=DB_DIR=/var/lib/nullsink
ExecStart=/opt/nullsink/deploy/backup.sh

# Hardening (mirrors the app units). NB: ProtectHome=read-only still lets the push step read an
# ssh/rclone credential under the service user's home; the push needs network, so no IPAddress* deny here.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
