# core/deploy/monero-wallet-rpc-watchdog.service
[Unit]
Description=Restart monero-wallet-rpc if it is hung (alive but not answering RPC)
# Best-effort ordering; the script itself guards on the wallet unit being enabled + active, so it is
# harmless if that unit is absent.
After=monero-wallet-rpc.service
# Fires ONLY on a genuine watchdog fault now — a script error, or a restart that didn't bring the wallet
# back. A SUCCESSFUL heal pages itself as a wallet recovery (via alert.sh) and exits 0, so this never
# double-counts a heal as a watchdog failure. Same sink as the other units (unit name only, no journal).
OnFailure=status-alert@%n.service

[Service]
Type=oneshot
# Runs as root (no User=): it must `systemctl restart monero-wallet-rpc`, which needs privilege. It touches
# only the localhost wallet RPC and systemctl — no app data, no secrets.
ExecStart=/opt/nullsink/deploy/wallet-rpc-watchdog.sh
# A wedged wallet may ignore SIGTERM, so its stop leg can run the full TimeoutStopSec before SIGKILL. Give
# ExecStart room for the probe (~55s) + that stop + the heal page so the watchdog is never itself SIGKILLed
# mid-heal (which would false-page a timeout). TimeoutStopSec on the wallet unit keeps that stop short.
TimeoutStartSec=180
