# core/deploy/monero-wallet-rpc-watchdog.timer
# Probe the wallet's liveness every few minutes and bounce it if wedged. setup.sh + deploy.sh enable this
# (via lib.sh enable_timers) whenever monero-wallet-rpc is enabled and disable it otherwise, so it tracks the
# XMR rail automatically. Manual: systemctl enable --now monero-wallet-rpc-watchdog.timer.
# Monotonic timers resume on their own after a reboot; Persistent= is intentionally omitted (it only affects
# OnCalendar= timers, and a missed liveness probe is worthless to replay late).
[Unit]
Description=Probe monero-wallet-rpc liveness (hung-wallet watchdog)

[Timer]
# A wedge strands the buy rail (deposits stop crediting), so check more often than the 10-min status-check:
# 3 min bounds the worst-case blind window. The script's retry loop covers a just-restarted wallet's reopen,
# so no startup-grace/interval coupling is needed.
OnBootSec=5min
OnUnitActiveSec=3min
AccuracySec=30s

[Install]
WantedBy=timers.target
