# core/deploy/nullsink-proxy.service
[Unit]
Description=nullsink proxy — metered LLM proxy (proxy trust domain)
After=network-online.target
Wants=network-online.target
# Page Telegram if this unit enters `failed` (a clean crash, or a crash-loop that trips the start limit).
# Secondary to status-check.timer, which catches a hung-but-"active" or stalled service every 10 min.
OnFailure=status-alert@%n.service

[Service]
Type=simple
# ONE service user for both app services, for now. Splitting uids would force splitting
# StateDirectory (systemd chowns it to the unit's User=), which moves balances.db / pending.db to new
# paths — and `nsk financials` opens BOTH, as a non-root user (cli/guard.ts refuses root: a root open
# strands root-owned WAL sidecars that break billing writes). So the uid split is blocked on the
# admin-plane redesign, and lands with it. Until then the socket's owner-only mode (below) is the whole
# gate, and it is sufficient: a same-uid attacker could write balances.db directly.
User=nullsink
Group=nullsink
WorkingDirectory=/opt/nullsink
# Shared with nullsink-payments.service. Each service reads only the vars it needs and ignores the rest
# (this one never reads MONERO_*/BITCOIN_*; payments never reads *_API_KEY). Splitting the file needs
# per-service ownership, so it lands with the uid split above.
#
# CAREFUL: systemd lets EnvironmentFile= OVERRIDE Environment=, not the reverse. Every Environment= line below
# therefore holds only because setup.sh's templated env file omits that variable. Don't add DB_PATH /
# CREDIT_SOCK / PENDING_DB_PATH to /etc/nullsink.env — they would silently win over the unit.
EnvironmentFile=/etc/nullsink.env
# Runs the self-contained `bun build --compile` artifact via the `current-proxy` symlink, swapped atomically
# by deploy.sh/setup.sh's install_binary (versioned binary + relative symlink for rollback).
ExecStart=/usr/local/lib/nullsink/current-proxy
# Type=simple considers ExecStart launched before the process has bound its owner-only credit socket. Gate the
# start job on that socket so nullsink-payments.service's After= ordering covers actual crossing readiness,
# not merely process creation. TimeoutStartSec bounds a broken/mismatched CREDIT_SOCK boot and lets systemd
# fail + restart the proxy instead of holding the payments start job forever.
ExecStartPost=/bin/sh -c 'until [ -S /run/nullsink/credit.sock ]; do sleep 0.1; done'
TimeoutStartSec=60
Restart=always
RestartSec=2

# Balance ledger lives here: /var/lib/nullsink, created service-owned and writable even under
# ProtectSystem=strict. This is the only stateful asset this service owns — back it up (losing it =
# refunding everyone). It holds token hashes + balances + the applied-credit markers; no identity, no
# request history. pending.db is in the same directory but belongs to nullsink-payments.service.
StateDirectory=nullsink
Environment=DB_PATH=/var/lib/nullsink/balances.db

# The credit crossing. THIS service binds the socket; nullsink-payments connects to it. RuntimeDirectory is
# created service-owned and torn down when this unit stops — payments then sees ENOENT, which its durable
# outbox treats as ambiguous (retry next tick), never as a lost credit. 0710 keeps the directory
# untraversable by anyone but the service user and its group; the socket itself is bound owner-only (0700)
# by src/credit-server.ts, and a connect(2) needs WRITE permission on that file — a kernel-enforced gate.
RuntimeDirectory=nullsink
RuntimeDirectoryMode=0710
Environment=CREDIT_SOCK=/run/nullsink/credit.sock

# Holds many concurrent, minutes-long upstream connections.
LimitNOFILE=1048576

# Memory + task containment. The box (4 GB) also runs Caddy, monero-wallet-rpc, Tor, and the OS, so a leak
# or flood in this service must not OOM the whole box. MemoryMax is a hard cgroup cap: exceeding it invokes
# the cgroup OOM killer, which reaps only processes in THIS unit (never Caddy/sshd box-wide), and
# OOMPolicy=stop + Restart=always relaunch cleanly. MemoryHigh applies reclaim back-pressure first (soft),
# below the hard cap. A post-OOM kill strands any in-flight holds, the same accepted edge as any ungraceful
# death (see proxy.ts shutdown notes). TasksMax bounds a thread/fork runaway. The two services' envelopes
# are deliberately asymmetric: this service streams 32 MiB bodies over long-lived connections and gets the
# bulk; payments handles tiny JSON and gets 256M. Their MemoryMax sums to 2048M, leaving roughly 2 GB for
# everything else on the box.
MemoryHigh=1408M
MemoryMax=1792M
TasksMax=512
OOMPolicy=stop

# Drain in-flight requests on stop/restart so their billing settles (the app handles SIGTERM: it lets
# streams finish for a short grace, then bills any straggler its metered partial + refunds the rest before
# force-closing). This window must stay ABOVE the app's SHUTDOWN_GRACE_MS (default 25s) so the app always
# finishes settlement on its own terms; if systemd SIGKILLs first, an in-flight hold is left un-reconciled.
KillSignal=SIGTERM
TimeoutStopSec=60

# Hardening. The app writes only the balance ledger (StateDirectory) and the credit socket
# (RuntimeDirectory); ProtectSystem=strict keeps the rest of the filesystem read-only.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
